# ooo.conf
[program:ooo]
command=python3 /home/curator2/ooo/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/ooo/test.err.log
stdout_logfile=/home/curator2/ooo/test.out.log